# jellytrack.service
[Unit]
Description=Jellytrack - Jellyfin Playback Tracker
After=network.target

[Service]
Type=simple
User=alexander-hinz
Group=alexander-hinz
WorkingDirectory=/home/alexander-hinz/jellytrack
Environment=PATH=/home/alexander-hinz/jellytrack/.venv/bin:/usr/bin
EnvironmentFile=/home/alexander-hinz/jellytrack/.env
ExecStart=/home/alexander-hinz/jellytrack/.venv/bin/python -m src.main
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jellytrack

[Install]
WantedBy=multi-user.target
