# jellytrack.service
[Unit]
Description=Jellytrack - Jellyfin Playback Tracker
After=network.target

[Service]
Type=simple
User=stinz
Group=stinz
WorkingDirectory=/home/stinz/jellytrack
Environment=PATH=/home/stinz/jellytrack/.venv/bin:/usr/bin
EnvironmentFile=/home/stinz/jellytrack/.env
ExecStart=/home/stinz/jellytrack/.venv/bin/python -m src.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
